# nexrad.service
# ── nexrad.service ────────────────────────────────────────────────────────────
# Systemd unit for the NEXRAD Weather Radar Display (gunicorn production mode).
#
# Install:
#   sudo cp nexrad.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable nexrad
#   sudo systemctl start nexrad
#
# Assumes the app lives at /opt/nexrad with a venv at /opt/nexrad/venv.
# Run as a dedicated non-root user:
#   sudo useradd -r -s /sbin/nologin -d /opt/nexrad nexrad
#   sudo chown -R nexrad:nexrad /opt/nexrad
# ─────────────────────────────────────────────────────────────────────────────

[Unit]
Description=NEXRAD Weather Radar Display
Documentation=https://github.com/YOUR_USERNAME/nexrad-radar
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nexrad
Group=nexrad
WorkingDirectory=/opt/nexrad

# Set NEXRAD_SECRET_KEY to a long random string in production.
# Generate one with: python3 -c "import secrets; print(secrets.token_hex(32))"
Environment="NEXRAD_SECRET_KEY=CHANGE_THIS_TO_A_RANDOM_STRING"
Environment="NEXRAD_DEBUG=0"

ExecStart=/opt/nexrad/venv/bin/gunicorn \
    --workers 2 \
    --bind 127.0.0.1:5000 \
    --timeout 30 \
    --keep-alive 5 \
    --access-logfile /opt/nexrad/logs/access.log \
    --error-logfile /opt/nexrad/logs/error.log \
    app:app

Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# ── Systemd security hardening ────────────────────────────────────────────────
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
# Allow writes only to the app's own data directories
ReadWritePaths=/opt/nexrad/logs /opt/nexrad/settings.json /opt/nexrad/logs/notified.json

[Install]
WantedBy=multi-user.target
