# wetnose-watchdog.service
# ── wetnose-watchdog.service ─────────────────────────────────────────────────
# One-shot health check: if /api/health doesn't respond, restart wetnose.service.
# Triggered on a fixed cadence by wetnose-watchdog.timer.
#
# This catches hangs and deadlocks that systemd's Restart=on-failure can't see
# (the Flask process is alive, but the event loop is wedged).
#
# Install:
#   sudo cp wetnose-watchdog.service wetnose-watchdog.timer /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now wetnose-watchdog.timer
# ─────────────────────────────────────────────────────────────────────────────

[Unit]
Description=Wet Nose Weather health-check watchdog
After=wetnose.service
Requisite=wetnose.service

[Service]
Type=oneshot
# Probe health. If unreachable, log it and bounce the service.
ExecStart=/bin/sh -c '\
  if ! /usr/bin/curl -fsS --max-time 5 http://127.0.0.1:5000/api/health >/dev/null; then \
    /usr/bin/logger -t wetnose-watchdog "Health check failed; restarting wetnose"; \
    /bin/systemctl restart wetnose.service; \
  fi'

# Hardening for the watchdog itself.
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
