# wetnose-watchdog.timer
# ── wetnose-watchdog.timer ───────────────────────────────────────────────────
# Fires wetnose-watchdog.service every 5 minutes after boot.
# ─────────────────────────────────────────────────────────────────────────────

[Unit]
Description=Run Wet Nose Weather watchdog every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=15s
Persistent=true
Unit=wetnose-watchdog.service

[Install]
WantedBy=timers.target
